# 40simta.cron
# simta admin tasks
10 00 * * * root    /usr/sbin/simtamaint
*/5 * * * * root    /usr/sbin/simc -q &>/dev/null
